# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
autostart=true
user=root
command=/usr/bin/redis-server

[program:lita]
user=root
directory=/srv/lita-docker
command=/usr/local/bin/bundle exec lita start
